# flaskApp.conf
[program:flaskApp]
user=root
directory=/home/ubuntu
command=python3 server.py
autostart=true
autorestart=false
stdout_logfile=/home/ubuntu/flaskApp.log
stderr_logfile=/home/ubuntu/flaskApp.err
